# config/crontab
* * * * * /rails/bin/rake crawler:crawl_all
